# data/scarecrow-session-monitor.service.in
[Unit]
Description=Monitor Session leader for Scarecrow Session
CollectMode=inactive-or-failed

# All services started after scarecrow-session-pre.target need to be torn down
# before the session finish can be signalled back to the display manager.
PartOf=scarecrow-session-pre.target
Before=scarecrow-session-pre.target

# No After, as we want this to start up immediately

[Service]
Type=notify
ExecStart=@libexecdir@/scarecrow-session-ctl --monitor
TimeoutStopSec=5
