# data/scarecrow-session-signal-init.service.in
[Unit]
Description=Signal initialization done to Scarecrow Session Manager

PartOf=scarecrow-session.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/scarecrow-session-ctl --signal-init
